# crontab
#run python script every minutes
* * * * * /usr/local/bin/python ./main.py
# An empty line is required at the end of this file for a valid cron file.
